# meta-titania/recipes-titania/http-splash/files/http-splash.socket
[Unit]
Description=Boot time HTTP socket
Before=nginx.service

[Socket]
ListenStream=80
Accept=yes

[Install]
WantedBy=sockets.target