# meta-titania/recipes-titania/datafs-resizer/files/datafs-resizer.service
[Unit]
Description=SD card size check
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/datafs-resizer.sh
