# meta-titania/recipes-titania/datafs-resizer/files/datafs-resizer.service
[Unit]
Description=SD card size check
Before=docker.service

[Service]
Restart=no
ExecStart=/usr/bin/datafs-resizer.sh
