# meta-titania/recipes-titania/datafs-resizer/files/var-lib-docker.automount
# TODO: maybe part of Docker bbappend?
[Unit]
Description=Auto-mount Data SD card volumne
Requires=datafs-resizer.service

[Automount]
Where=/var/lib/docker

[Install]
WantedBy=multi-user.target
