# meta-titania/recipes-support/swupdate/files/after-everything.target
[Unit]
Description=Target for units that must be started the last
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes