# meta-titania/recipes-support/swupdate/files/check-update.service
[Unit]
Description=Check for successful update
After=multi-user.target

[Service]
ExecStart=/sbin/check_update.sh
Type=simple

[Install]
WantedBy=after-everything.target
